# pkg/centos/after-install.sh
/sbin/chkconfig --add keyless

chown -R keyless:keyless /etc/keyless
