# pkg/centos/after-install.sh
/sbin/chkconfig --add gokeyless

chown -R keyless:keyless /etc/keyless
